# src/etc/systemd/system/powerwatcher.service
[Unit]
Description="Powerwatcher service"

[Service]
ExecStart=/usr/local/bin/powerwatcher.sh
